# tetsuya.service
[Unit]
Description=Tetsuya server (auto-upgrade via uvx)
After=network-online.target
Wants=network-online.target

[Service]
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env uvx tetsuya@latest --logistro-level INFO server start
Restart=on-failure
RestartSec=5
StartLimitInterval=300
StartLimitBurst=5

[Install]
WantedBy=default.target
